# deployment/templates/pandora-files.service.j2
[Unit]
Description=Pandora files backend
After=network.target

[Service]
Type=simple
User={{ user }}
ExecStart=node {{ dest_dir }}/app.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
